# docker/scripts/supervisor/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D -e
stdout_events_enabled = true
stderr_events_enable = true
autorestart=true
auto_start = true
