# docker/scripts/supervisor/uwsgi.conf
[program:#NAME#]
command = /usr/bin/uwsgi --ini /data/scripts/uwsgi.ini
stdout_events_enabled = true
stderr_events_enable = true
auto_start = true
autorestart = true
